# supervisord.conf
; supervisor config file

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700                  ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisord.log ; (main log file; default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile; default $CWD/supervisord.pid)
childlogdir=/var/log/supervisor  ; ('AUTO' child log dir, default $TEMP)

; The section below must be in the present for the RPC (supervisorctl/web)
; interface in to function.
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

[program:dataslicer]
;command=/usr/local/bin/bokeh serve myapp.py --host foo.com:80
command=/usr/local/bin/bokeh serve analyse.py --port 5100 --host 188.166.74.162:5100 --host dataslicer.qualikiz.com:80
directory=/home/bokeh/QuaLiKiz-dataslicer
autorestart=true
startretries=0
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/%(program_name)s.err.log
stdout_logfile=/var/log/%(program_name)s.out.log
user=bokeh
environment=USER="bokeh",HOME="/home/bokeh"
